# service/leroy.service
[Unit]
Description=Leroy
After=multi-user.target

[Service]
Type=idle
WorkingDirectory=/home/pi/Projects/project-leroy
ExecStart=/home/pi/Projects/project-leroy/run.sh
Restart=on-abort
User=pi
Environment=DISPLAY=:0

[Install]
WantedBy=multi-user.target